# app/launchd/com.timeclaw.radio-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.timeclaw.radio-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/Cellar/node@22/22.22.0/bin/node</string>
    <string>/Users/doorliss/timeclaw/radio-web/sync.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/doorliss/timeclaw/radio-web</string>

  <key>StartInterval</key>
  <integer>900</integer> <!-- 15 minutes -->

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/doorliss/timeclaw/radio-web/logs/sync.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/doorliss/timeclaw/radio-web/logs/sync.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/local/bin</string>
  </dict>
</dict>
</plist>
