# watcher/launchd/com.curtis.whatsapp-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.curtis.whatsapp-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/doorliss/.openclaw/workspace/whatsapp-watcher/watcher.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/doorliss/.openclaw/workspace/whatsapp-watcher</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/doorliss</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <!-- Throttle restarts so a crashloop won't spin a CPU. -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/doorliss/.openclaw/workspace/whatsapp-watcher/logs/watcher.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/doorliss/.openclaw/workspace/whatsapp-watcher/logs/watcher.err.log</string>
</dict>
</plist>
