# watcher/launchd/com.curtis.whatsapp-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.curtis.whatsapp-digest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/doorliss/.openclaw/workspace/whatsapp-watcher/digest.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/doorliss/.openclaw/workspace/whatsapp-watcher</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/doorliss</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <!-- 4 digests per day at 8am, 1pm, 6pm, 11pm local time, each covering the
       6 hours since the last one. Tuned to message rhythm: 9am spike caught at
       1pm, 7pm spike caught at 11pm, etc. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/doorliss/.openclaw/workspace/whatsapp-watcher/logs/digest.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/doorliss/.openclaw/workspace/whatsapp-watcher/logs/digest.err.log</string>
</dict>
</plist>
